# scripts/ipinfo/rc.local
#!/bin/sh

# Jonas: these lines write output to specified logfile
exec 2> /path/to/tmp/rc.local.log      # send stderr from rc.local to a log file
exec 1>&2                              # send stdout to the same log file
set -x                                 # tell sh to display commands before execution
# Jonas: this script sends an Email containing the current local ip for easy ssh access.
python /path/to/scripts/ipinfo/ip_mailer.py

exit 0 # needed
